# scrapers/persrm-scraper.service
[Unit]
Description=PersRM Market Data Scraper Service
After=network.target

[Service]
Type=simple
User=kr
Group=kr
WorkingDirectory=/home/kr/ChatOS-v2.0/scrapers
ExecStart=/home/kr/text-generation-webui/venv/bin/python /home/kr/ChatOS-v2.0/scrapers/market_scraper.py
Restart=always
RestartSec=10
StandardOutput=append:/home/kr/ChatOS-v2.0/scrapers/logs/service.log
StandardError=append:/home/kr/ChatOS-v2.0/scrapers/logs/service.log

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/kr

[Install]
WantedBy=multi-user.target

